# deploy/snek-server.service
# snek3 trainer-server as a systemd service (for a VM / bare-metal GPU box you
# own — NOT needed on managed GPU pods, which are already containers).
#
# Install:
#   sudo cp deploy/snek-server.service /etc/systemd/system/
#   # edit User=, WorkingDirectory=, the token, and any training flags below
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now snek-server
#   journalctl -u snek-server -f          # logs (incl. the dashboard URL)
#
# Starts an IDLE server (create runs from the dashboard). To auto-start a run,
# append e.g. --run-id myrun to ExecStart.

[Unit]
Description=snek3 training server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brensch
WorkingDirectory=/home/brensch/snek3
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True
Environment=SNEK_SERVE_TOKEN=CHANGEME
# Direct python (no make) so a restart never rebuilds. Mirrors `make server`
# defaults so dashboard-created runs share the meta.
ExecStart=/home/brensch/snek3/.venv/bin/python -m azsnek.train_albatross \
    --serve-host 0.0.0.0 --serve-port 8050 \
    --num-snakes 2 --count 64 --samples 8000 --depth 2 --iters 120 \
    --tau-min 0.5 --tau-max 10.0 --response-tau 12.0 --response-after 0 \
    --draw-value -0.9 --exploration-prob 0.15 --max-turns 0 --eval-batch-size 2048 \
    --filters 64 --blocks 6 --lr 1e-3 --train-steps 128 --recency 2.0 \
    --batch-size 1024 --buffer-size 500000 --eval-every 5 --eval-games 64 \
    --record-games 2 --record-every 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
